# deploy/supervisord.conf
[supervisord]
nodaemon=false

[program:uwsgi]
command = /usr/bin/uwsgi --ini  /etc/uwsgi.ini

[program:nginx]
command = /usr/sbin/nginx -g "daemon off;"
